# deployments/local/conf.d/supervisor/pm2.conf
[program:pm2d]
command = pm2-runtime /opt/docker/etc/pm2/ecosystem.config.js
process_name=%(program_name)s
startsecs = 0
autostart = false
autorestart = true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
